# .docker/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/gunicorn -b 0.0.0.0:8001 tpp.b24online_wsgi
directory=/src/b24project
stdout_logfile=/src/log/stdout.log
stderr_logfile=/src/log/stderr.log